# Shell/com.imperium.deskflow-client.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Headless Deskflow client (com.imperium.deskflow-client).

  Template — __HOME__ is substituted by token-api's _ensure_deskflow_client_agent()
  when syncing to ~/Library/Application Support/Imperium/. Its own launchd job, never
  a token-api child: token-api's job has no AbandonProcessGroup, so children die on
  every deploy restart (Mac analog of the WSL transient-unit fix, PR #148).

  RunAtLoad=false + KeepAlive=false: one-shot per invite. If the supervisor goes
  quiet (bounded connect/reconnect window expired), only a WSL satellite invite or a
  manual `launchctl kickstart` revives it — by design; the satellite owns retry.

  No EnvironmentVariables block on purpose: $XDG_CONFIG_HOME would override the
  core's ~/Library/Deskflow/Deskflow.conf resolution.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.imperium.deskflow-client</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__HOME__/Library/Application Support/Imperium/deskflow-client-supervisor.py</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Application Support/Imperium/deskflow-client.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Application Support/Imperium/deskflow-client.log</string>
</dict>
</plist>
