# cli-tools/launchd/ai.openclaw.process-gc-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ai.openclaw.process-gc-watchdog — Workstream C of the tmux lag remediation.
  Interval-driven (StartInterval, no KeepAlive): every 5 min it reaps orphaned
  (ppid=1) search processes that are egregiously stuck (high CPU + long-lived).
  Tunables below are env overrides for the script's defaults — edit here to tune
  without touching code. GC_PATTERN may be added here to opt build tools in (see
  the script header for the daemon-kill caveat).
  Log growth is bounded by in-script self-truncation: the script keeps the last
  ~1000 lines of its own ~/.claude/process-gc-watchdog.log once it exceeds 2000
  (one summary line per run); the StdOut/StdErr logs below stay ~empty.
  ProgramArguments points at a LOCAL copy of the script, NOT the NAS path. launchd's
  context cannot read/execve a file off the SMB mount (macOS Sequoia TCC gates
  network-volume access per-binary; system bash under launchd lacks it → exit 126
  "Operation not permitted", confirmed for both direct .sh exec and `/bin/bash <nas
  path>`). This mirrors the sibling tokenapi-watchdog, which also runs from
  ~/.local/bin. The script touches the NAS not at all at runtime (only ps/kill,
  ~/.claude/ logs, and localhost curl), so the local copy is fully self-contained.
  The repo copy at cli-tools/Shell/process-gc-watchdog.sh is canonical (version
  control); the local copy is the deployed runtime and must be refreshed from it on
  any change. Install/redeploy:
    cp <repo>/cli-tools/Shell/process-gc-watchdog.sh ~/.local/bin/process-gc-watchdog
    chmod +x ~/.local/bin/process-gc-watchdog
    cp <repo>/cli-tools/launchd/ai.openclaw.process-gc-watchdog.plist ~/Library/LaunchAgents/
    launchctl bootout  gui/$UID ~/Library/LaunchAgents/ai.openclaw.process-gc-watchdog.plist 2>/dev/null || true
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/ai.openclaw.process-gc-watchdog.plist
  (Use bootout+bootstrap, NOT kickstart -k — kickstart reuses the old plist.)
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>ai.openclaw.process-gc-watchdog</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/tokenclaw/.local/bin/process-gc-watchdog</string>
	</array>
	<key>StartInterval</key>
	<integer>300</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/tokenclaw/runtimes/Token-OS/live/cli-tools/bin:/Users/tokenclaw/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>/Users/tokenclaw</string>
		<key>GC_CPU_MIN</key>
		<string>40</string>
		<key>GC_ETIME_MIN</key>
		<string>30</string>
		<key>GC_NAS_PATH</key>
		<string>/Volumes/Imperium</string>
		<key>GC_NAS_ETIME_MIN</key>
		<string>10</string>
		<key>GC_ANCIENT_MIN</key>
		<string>60</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/tokenclaw/.claude/process-gc-watchdog-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/tokenclaw/.claude/process-gc-watchdog-stderr.log</string>
</dict>
</plist>
