# cli-tools/launchd/ai.openclaw.tokenapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ai.openclaw.tokenapi — LaunchAgent for the Mac token-api server (port 7777).

  Canonical, repo-controlled template. The LIVE plist at
  ~/Library/LaunchAgents/ai.openclaw.tokenapi.plist used to be a hand-edited,
  live-only artifact; it is now reproducible from here.

  Socket activation (the `Sockets` dict below) is the load-bearing change:
  launchd creates and OWNS the listening socket and hands its fd to the job via
  launch_activate_socket(3) (see token-api/launchd_socket.py + main.py's
  __main__ shim). Because launchd — not uvicorn — owns the socket, the kernel
  accept backlog survives every uvicorn restart: during a deploy, incoming hooks
  (most critically the SessionStart registration POST) stall briefly in the
  kernel instead of being connection-refused and silently dropped.

  ProgramArguments runs `python main.py` (NOT `-m uvicorn …`) so main.py's shim
  retrieves the activated fd and serves on it (host/port fallback otherwise).

  Install / refresh on an existing machine (idempotent, preserves secrets):
    token-restart            # reconciles ProgramArguments + Sockets in place
  then, because Sockets/ProgramArguments changes need a full reload (kickstart
  alone reuses the bootstrapped definition):
    launchctl bootout  gui/$(id -u)/ai.openclaw.tokenapi 2>/dev/null
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/ai.openclaw.tokenapi.plist

  Fresh install: copy this file to ~/Library/LaunchAgents/, then fill the
  EnvironmentVariables secrets below (the in-place reconciler never touches
  EnvironmentVariables, so live secrets are preserved across upgrades).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>ai.openclaw.tokenapi</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/tokenclaw/.local/venvs/token-api/bin/python</string>
		<string>main.py</string>
	</array>
	<key>Sockets</key>
	<dict>
		<key>Listeners</key>
		<dict>
			<key>SockServiceName</key>
			<string>7777</string>
			<key>SockType</key>
			<string>stream</string>
			<key>SockFamily</key>
			<string>IPv4</string>
		</dict>
	</dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/tokenclaw</string>
		<key>PATH</key>
		<string>/Users/tokenclaw/runtimes/Token-OS/live/cli-tools/bin:/Users/tokenclaw/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>BRAVE_API_KEY</key>
		<string>__SET_BY_OPERATOR__</string>
		<key>CD_RESTART_SECRET</key>
		<string>__SET_BY_OPERATOR__</string>
		<key>TOKEN_SATELLITE_REFRESH_SECRET</key>
		<string>__SET_BY_OPERATOR__</string>
		<key>TOKEN_API_OPENAI_TTS_INSTRUCTIONS</key>
		<string>Read briskly and clearly. Use a fast default speaking pace. Do not add theatrical pauses, sighs, breaths, gasps, inhalations, exhalations, or other nonverbal breathing sounds. Deliver the text directly without performative affect.</string>
		<key>TOKEN_API_DISABLE_AUTOMATIC_IMPERIUM_DAILY_NOTE_WRITES</key>
		<string>1</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>16384</integer>
	</dict>
	<key>HardResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>32768</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/tokenclaw/runtimes/Token-OS/live/token-api</string>
	<key>StandardOutPath</key>
	<string>/Users/tokenclaw/.claude/token-api-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/tokenclaw/.claude/token-api-stderr.log</string>
</dict>
</plist>
