# cli-tools/launchd/ai.tokenclaw.tmuxctld.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ai.tokenclaw.tmuxctld — LaunchAgent for the tmuxctld HTTP-loopback daemon.

  Install / refresh:
    cli-tools/bin/tmuxctld-ctl install
  (copies this plist to ~/Library/LaunchAgents/ then bootout + bootstrap)

  Or manually:
    cp cli-tools/launchd/ai.tokenclaw.tmuxctld.plist ~/Library/LaunchAgents/
    launchctl bootout  gui/$(id -u)/ai.tokenclaw.tmuxctld 2>/dev/null
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/ai.tokenclaw.tmuxctld.plist

  ProgramArguments point at the LIVE runtime checkout, so daemon code updates land
  via the normal live-checkout sync — there is no venv to sync (stdlib-only). Only
  the plist itself is synced to ~/Library/LaunchAgents/ (via `tmuxctld-ctl install`).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>ai.tokenclaw.tmuxctld</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/tokenclaw/runtimes/Token-OS/live/cli-tools/bin/tmuxctld</string>
		<string>--host</string>
		<string>127.0.0.1</string>
		<string>--port</string>
		<string>7778</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/tokenclaw</string>
		<key>PATH</key>
		<string>/Users/tokenclaw/runtimes/Token-OS/live/cli-tools/bin:/Users/tokenclaw/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>ProcessType</key>
	<string>Background</string>
	<key>WorkingDirectory</key>
	<string>/Users/tokenclaw/runtimes/Token-OS/live/cli-tools</string>
	<key>StandardOutPath</key>
	<string>/Users/tokenclaw/.claude/tmuxctld-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/tokenclaw/.claude/tmuxctld-stderr.log</string>
</dict>
</plist>
