# edge_proxy/systemd/edge-proxy.service
[Unit]
Description=Token-OS edge_proxy stateless routing proxy
After=network-online.target token-api.socket
Wants=network-online.target
# The proxy binds the box's tailnet IP; at boot that IP can lag tailscaled, so a
# failed bind must retry until the address exists instead of tripping the
# default start-limit (5 in 10s) and staying down.
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=%h/runtimes/Token-OS/live/edge_proxy
Environment=EDGE_PROXY_CONFIG=%h/secrets/token-os/edge_proxy.json
Environment=IMPERIUM_MACHINE=k12-personal
ExecStart=%h/.bun/bin/bun src/server.ts
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=default.target
