# k12_daemon/systemd/k12-daemon.service
[Unit]
Description=Token-OS k12_daemon event-sourced tmux control plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/runtimes/Token-OS/live/k12_daemon
Environment=K12_DAEMON_CONFIG=%h/secrets/token-os/k12_daemon.json
Environment=IMPERIUM_MACHINE=k12-personal
ExecStart=%h/.bun/bin/bun src/daemon.ts
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
# Terminate only the bun main process on stop/restart — NOT the whole cgroup.
# The default (control-group) SIGTERMs the daemon-spawned tmux server, killing
# the persistent estate on every `systemctl restart`. The daemon otherwise
# spawns only short-lived tmux CLI children, so process-scope teardown is safe
# and leaves the estate (and its seat sessions) alive for clean idempotency.
KillMode=process

[Install]
WantedBy=default.target
